# deploy/systemd/hfinance.service
[Unit]
Description=HFinance - 个人资产管理记账系统
After=network.target

[Service]
Type=simple
User=__DEPLOY_USER__
WorkingDirectory=__PROJECT_DIR__/backend
ExecStart=__PROJECT_DIR__/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8820
Restart=always
RestartSec=5
EnvironmentFile=/etc/default/hfinance

[Install]
WantedBy=multi-user.target
